# deploy/kars.service
[Unit]
Description=KARS Media Archive
After=network.target

[Service]
Type=simple
User=kars
Group=kars
WorkingDirectory=/opt/kars
ExecStart=/opt/kars/kars
EnvironmentFile=/opt/kars/.env
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
